# deploy/chat-analysis.service
[Unit]
Description=Gunicorn instance for Chat Analysis
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/chat-analysis
Environment="PATH=/home/ubuntu/chat-analysis/venv/bin"
ExecStart=/home/ubuntu/chat-analysis/venv/bin/gunicorn --workers 1 --timeout 120 --bind 127.0.0.1:5328 run:app

[Install]
WantedBy=multi-user.target
